# dot_config/systemd/user/mihomo.service
[Unit]
Description=Mihomo local proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/bin/test -s %h/.config/mihomo/config.yaml
ExecStartPre=/usr/bin/mihomo -t -d %h/.local/share/mihomo -f %h/.config/mihomo/config.yaml
ExecStart=/usr/bin/mihomo -d %h/.local/share/mihomo -f %h/.config/mihomo/config.yaml
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
