# dot_config/systemd/user/vicinae.service.d/ime.conf
[Service]
# Desktop entries launched by Vicinae must inherit the IME bridge variables.
Environment=XMODIFIERS=@im=fcitx
Environment=GTK_IM_MODULE=fcitx
Environment=QT_IM_MODULE=fcitx
Environment="QT_IM_MODULES=wayland;fcitx;ibus"
Environment=SDL_IM_MODULE=fcitx
Environment=GLFW_IM_MODULE=ibus
